# host/systemd/fwrouter-traffic-collect.timer
[Unit]
Description=Run FWRouter v2 traffic accounting collect every 3 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=3min
AccuracySec=15s
Persistent=true
Unit=fwrouter-traffic-collect.service

[Install]
WantedBy=timers.target
