# gateway-stack/fwrouter-stack/host-systemd/fwrouter-apply-devices.service
[Unit]
Description=Apply fwrouter device overrides (ipsets only)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fwrouter-apply-devices
