# gateway-stack/fwrouter-stack/host-systemd/fwrouter-apply.service
[Unit]
Description=Apply fwrouter firewall/ipset rules
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fwrouter-apply --apply

[Install]
WantedBy=multi-user.target
