# fwrouter-stack/host-systemd/dnsmasq.service.d/10-ipset-cap.conf
[Service]
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
ExecStartPre=/usr/local/sbin/fwrouter-ipset-init
