# fwrouter-stack/host-systemd/fwrouter-sync-tailscale.path
[Unit]
Description=Watch tailscale state for peer changes

[Path]
PathModified=/var/lib/tailscale/tailscaled.state
Unit=fwrouter-sync-tailscale.service

[Install]
WantedBy=multi-user.target
