# fwrouter-stack/host-systemd/fwrouter-apply-devices.path
[Unit]
Description=Watch fwrouter devices.conf for changes

[Path]
PathModified=/etc/fwrouter/devices.conf
Unit=fwrouter-apply-devices.service

[Install]
WantedBy=multi-user.target
