# fwrouter-stack/host-systemd/fwrouter-health-check.service
[Unit]
Description=FWRouter periodic health check and self-heal
After=network-online.target docker.service tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fwrouter-health-check

