# fwrouter-stack/host-systemd/fwrouter-sync-tailscale.service
[Unit]
Description=Sync tailscale peers into fwrouter devices
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fwrouter-sync-tailscale
