# fwrouter-stack/host-systemd/fwrouter-resolve-domains.path
[Unit]
Description=Watch fwrouter rules for changes

[Path]
PathModified=/etc/fwrouter/rules.conf
PathModified=/etc/fwrouter/rules.d
Unit=fwrouter-resolve-domains.service

[Install]
WantedBy=multi-user.target
