# fwrouter-stack/host-systemd/fwrouter-apply.service.d/10-deps.conf
[Unit]
After=network-online.target docker.service netfilter-persistent.service
Wants=network-online.target netfilter-persistent.service
