# fwrouter-stack/host-systemd/netfilter-persistent.service.d/10-fwrouter-ipset.conf
[Service]
ExecStartPre=/usr/local/sbin/fwrouter-ipset-init
